# telegram-bot.service
[Unit]
Description=Telegram Bot with Claude AI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bot
Environment="PATH=/opt/bot/venv/bin"
ExecStart=/opt/bot/venv/bin/python /opt/bot/run_bot.py
StandardOutput=append:/opt/bot/bot.log
StandardError=append:/opt/bot/bot.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target



